# ai-email-poller.service.example
[Unit]
Description=AI Email Filter Poller
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-email-filter
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/ai-email-filter/venv/bin/python -m services.poller
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target